# site-files/gunicorn.conf
description "Gunicorn daemon for Factitious project"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [!12345]

# If the process quits unexpectadly trigger a respawn
respawn

setuid django
setgid django
chdir /home/django

exec gunicorn \
    --name=factitious_project \
    --pythonpath=factitious_project/server/factitious_app \
    --bind=127.0.0.1:9000 \
    --config /etc/gunicorn.d/gunicorn.py \
    django_project.wsgi:application
